# docker-supervisord.conf
[supervisord]
nodaemon=true

[program:bn]
command=node server/app.js -i 0 --no-daemon

[program:fn]
command=http-server
